# profile.d/tldr.sh
# · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · ·
# Add settings for tldr (https://github.com/tldr-pages/tldr)
#
# Author: oqkr@pm.me
# Source: https://github.com/oqkr/dotfiles
# · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · ·

command -v tldr >/dev/null || return

alias tldr='tldr --theme ocean'
