# profile.d/vim.sh
# · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · ·
# Add settings for Vim.
#
# Author: oqkr@pm.me
# Source: https://github.com/oqkr/dotfiles
# · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · ·

command -v vim >/dev/null || return

alias vi="vim"
export EDITOR="vim"
export VISUAL="vim"
