# profile.d/z.sh
# · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · ·
# Enable z (https://github.com/rupa/z).
#
# Author: oqkr@pm.me
# Source: https://github.com/oqkr/dotfiles
# · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · · ·

# This is the default path when installed by homebrew.
[[ -f /usr/local/etc/profile.d/z.sh ]] || return

source /usr/local/etc/profile.d/z.sh
